# meta-uhmi/meta-distributed-display-fw/recipes-core/ucl-tools/files/ucl-launcher.service
[Unit]
Description=ucl-launcher

[Service]
Type=idle
ExecStart=/usr/bin/ucl-launcher
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
